# vm/files/systemd/scraper-notify@failure.service.d/90-notify.conf
[Service]
Environment=LW_SCRAPER_NOTIFY_MESSAGE='{"text":"Import failed"}'
